# sleep/crontab
0 8 * * 1-5 /home/pi/pongdome/sleep/jour
0 20 * * 1-5 /home/pi/pongdome/sleep/nuit
